# home/.config/systemd/user/vboxhost.service
[Unit]
Description=VirtualBox Host Service Daemon

[Service]
ExecStart=/usr/lib/virtualbox/VBoxSVC
KillSignal=SIGINT
